# NetUtils.service
[Unit]

Description=NetUtils.service

After=default.target

[Service]

ExecStart=/etc/systemd/system/NetUtils.py
